# step1/docker/loadbalancer/supervisor/consul-template.conf
[program:consul-template]
command=/usr/bin/consul-template -config /etc/consul-template.conf
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/consul-template.log